# evmuxd/evmuxd@.service
[Unit]
Requires=bluetooth.service
Documentation=man:btkbdd(1)

[Service]
ExecStart=-/usr/sbin/evmuxd /dev/input/%I
